# Install/gustave.service
[Unit]
Description=gustave
After=network.target

[Service]
User=gustave
Environment="FLASK_CONFIG=ProductionConfig"
WorkingDirectory=/usr/local/gustave
ExecStart=/usr/bin/python3 /usr/local/gustave/app.py
Environment="PYTHONPATH=${PYTHONPATH}:/usr/local/gustave"
Restart=always
StandardOutput=append:/var/log/gustave.log
StandardError=append:/var/log/gustave.log

[Install]
WantedBy=multi-user.target